# etc/systemd/system/bluetooth-agent.service
[Unit]
Description=Bluetooth Agent

[Service]
User=pi
Type=simple
ExecStart=/home/pi/bluetooth/bt-agent.py

[Install]
WantedBy=multi-user.target
