# etc/systemd/system/bmw-agent.service
[Unit]
Description = BMW Agent

[Service]
User=pi
Type=simple
ExecStart=/usr/bin/python /home/pi/bmw/agent.py

[Install]
WantedBy=multi-user.target
