# etc/systemd/system/bmw_agent.service
[Unit]
Description = BMW Agent

[Service]
User=pi
Type=simple
ExecStart=/home/pi/agent.py

[Install]
WantedBy=multi-user.target
